# etc/systemd/system/linkprofiler-coordinator.service
[Unit]
Description=Link Profiler Job Coordinator
After=network.target redis.service postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Link_Profiler_Repo
Environment=PATH=/opt/Link_Profiler_Repo/venv/bin
Environment=PYTHONPATH=/opt/Link_Profiler_Repo
Environment=LP_DATABASE_URL=postgresql://linkprofiler:secure_password_123@localhost:5432/link_profiler_db
Environment=DATABASE_URL=postgresql://linkprofiler:secure_password_123@localhost:5432/link_profiler_db
Environment=LP_REDIS_URL=redis://:redis_secure_pass_456@127.0.0.1:6379/0
Environment=REDIS_URL=redis://:redis_secure_pass_456@127.0.0.1:6379/0
Environment=ACCESS_TOKEN=your_secure_api_token_789
Environment=LP_ACCESS_TOKEN=your_secure_api_token_789
Environment=LP_AUTH_SECRET_KEY=cb53bd28ffc37f3f99bcae7b2f23252fc74fe233f623c37b4fbce541c0c672d3
ExecStart=/opt/Link_Profiler_Repo/venv/bin/python Link_Profiler/queue_system/job_coordinator.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
