# etc/systemd/system/linkprofiler-monitoring.service
[Unit]
Description=Link Profiler Monitoring Dashboard
After=network.target linkprofiler-api.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/Link_Profiler_Repo/Link_Profiler/monitoring
Environment=PATH=/opt/Link_Profiler_Repo/venv/bin
Environment=PYTHONPATH=/opt/Link_Profiler_Repo
Environment=LP_DATABASE_URL="postgresql://linkprofiler:secure_password_123@localhost:5432/link_profiler_db"
Environment=DATABASE_URL=postgresql://linkprofiler:secure_password_123@localhost:5432/link_profiler_db
Environment=LP_REDIS_URL=redis://:redis_secure_pass_456@127.0.0.1:6379/0
Environment=REDIS_URL=redis://:redis_secure_pass_456@127.0.0.1:6379/0
ExecStart=/opt/Link_Profiler_Repo/venv/bin/python /opt/Link_Profiler_Repo/Link_Profiler/monitoring/dashboard.py dashboard
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
